# systemd/llama-metrics.service
[Unit]
Description=LlamaMetrics dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=llama-metrics
Group=llama-metrics
WorkingDirectory=/opt/llama-metrics
Environment=HOME=/var/lib/llama-metrics
ExecStart=/opt/llama-metrics/.venv/bin/python -m app.run
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
